# services/precotracker-bot.service
[Unit]
Description=PriceTracker - Telegram Two-Way Bot
After=network.target

[Service]
Type=simple
User=iago
WorkingDirectory=/home/iago/Documents/Pessoal/price-tracker/backend
ExecStart=/home/iago/Documents/Pessoal/price-tracker/.venv/bin/python telegram_bot.py
Restart=on-failure

[Install]
WantedBy=multi-user.target
